# deploy/systemd/tron.service
[Unit]
Description=Tron Agent Server
Documentation=https://github.com/your-org/tron
After=network.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__

# Main executable
ExecStart=/usr/bin/node __TRON_INSTALL_PATH__/packages/server/dist/index.js

# Working directory
WorkingDirectory=__TRON_INSTALL_PATH__

# Environment
Environment=NODE_ENV=production
Environment=HOME=__HOME__
Environment=TRON_DATA_DIR=__HOME__/.tron
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Restart configuration
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=10

# Resource limits
LimitNOFILE=8192
LimitNPROC=4096

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__HOME__/.tron
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tron

# Timeout configuration
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
