# packages/mac-app/Sources/Resources/Library/LaunchAgents/com.tron.server.dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tron.server.dev</string>

    <key>ProgramArguments</key>
    <array>
        <string>tron</string>
        <string>--port</string>
        <string>9848</string>
        <string>--quiet</string>
    </array>

    <key>BundleProgram</key>
    <string>Contents/Library/LoginItems/Tron Server.app/Contents/MacOS/tron</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>TRON_HOME_NAME</key>
        <string>.tron-dev</string>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.tron.mac</string>
        <string>com.tron.mac.dev</string>
    </array>
</dict>
</plist>
